# internal/contain/images/agent/bootstrap/services/bootstrap-jobs/run
#!/bin/bash
# One-shot service: bootstraps default jobs on first container launch.
# After bootstrapping, removes its persisted runit service so it does not linger.
# Jobs are user-owned after creation: edit/delete via term-llm jobs CLI.

if [ "$(id -un)" != "agent" ]; then
  exec sudo -Hu agent env AGENT_NAME="${AGENT_NAME:-agent}" AGENT_DISTRO="${AGENT_DISTRO:-arch}" "$0" "$@"
fi

SENTINEL="/home/agent/.config/term-llm/.jobs-bootstrapped"
PERSISTED_SERVICE_DIR="/home/agent/.config/term-llm/services/bootstrap-jobs"
RUNSVDIR_LINK="/etc/runit/runsvdir/bootstrap-jobs"
SV_SERVICE_DIR="/etc/sv/bootstrap-jobs"
AGENT="${AGENT_NAME:-agent}"

cleanup_and_exit() {
  echo "bootstrap-jobs: removing one-shot service"
  rm -rf "$PERSISTED_SERVICE_DIR"
  sudo rm -f "$RUNSVDIR_LINK" 2>/dev/null || true
  sudo rm -rf "$SV_SERVICE_DIR" 2>/dev/null || true
  exit 0
}

# Already bootstrapped — remove this one-shot service if an older volume kept it.
if [ -f "$SENTINEL" ]; then
  cleanup_and_exit
fi

# Wait for jobs API to be ready (starts after runit launches the jobs service)
echo "bootstrap-jobs: waiting for jobs API..."
until nc -z 127.0.0.1 8080 > /dev/null 2>&1; do
  sleep 2
done
echo "bootstrap-jobs: jobs API ready, creating default jobs for agent=$AGENT"

# mine-sessions — extract memory fragments from session transcripts
term-llm jobs create --data "{
  \"name\": \"mine-sessions\",
  \"enabled\": true,
  \"runner_type\": \"program\",
  \"runner_config\": {
    \"command\": \"/usr/local/bin/term-llm\",
    \"args\": [\"memory\", \"mine\", \"--agent\", \"$AGENT\", \"--stats\", \"--limit\", \"50\", \"--insights\"]
  },
  \"trigger_type\": \"cron\",
  \"trigger_config\": { \"expression\": \"*/30 * * * *\", \"timezone\": \"UTC\" },
  \"concurrency_policy\": \"forbid\",
  \"timeout_seconds\": 3600,
  \"misfire_policy\": \"skip\"
}"

# update-recent — promote recent fragments into the agent's recent.md
term-llm jobs create --data "{
  \"name\": \"update-recent\",
  \"enabled\": true,
  \"runner_type\": \"program\",
  \"runner_config\": {
    \"command\": \"/usr/local/bin/term-llm\",
    \"args\": [\"memory\", \"update-recent\", \"--agent\", \"$AGENT\"]
  },
  \"trigger_type\": \"cron\",
  \"trigger_config\": { \"expression\": \"*/10 * * * *\", \"timezone\": \"UTC\" },
  \"concurrency_policy\": \"forbid\",
  \"timeout_seconds\": 300,
  \"misfire_policy\": \"skip\"
}"

SYSTEM_UPGRADE_ARGS='["pacman", "-Syu", "--noconfirm"]'
case "${AGENT_DISTRO:-arch}" in
  arch)
    SYSTEM_UPGRADE_ARGS='["pacman", "-Syu", "--noconfirm"]'
    ;;
  fedora)
    SYSTEM_UPGRADE_ARGS='["dnf", "-y", "upgrade"]'
    ;;
  *)
    echo "bootstrap-jobs: unknown AGENT_DISTRO=${AGENT_DISTRO}; defaulting system-upgrade to pacman" >&2
    ;;
esac

# system-upgrade — keep the container distro packages current
term-llm jobs create --data "{
  \"name\": \"system-upgrade\",
  \"enabled\": true,
  \"runner_type\": \"program\",
  \"runner_config\": {
    \"command\": \"sudo\",
    \"args\": $SYSTEM_UPGRADE_ARGS
  },
  \"trigger_type\": \"cron\",
  \"trigger_config\": { \"expression\": \"0 5 * * *\", \"timezone\": \"UTC\" },
  \"concurrency_policy\": \"forbid\",
  \"timeout_seconds\": 7200,
  \"misfire_policy\": \"skip\"
}"

touch "$SENTINEL"
echo "bootstrap-jobs: done — 3 jobs created. Edit them with: term-llm jobs list / term-llm jobs update"
cleanup_and_exit
