# docker/templates/services/memory-mine/run
#!/bin/sh
# Mine session transcripts into memory fragments every 6 hours.
# Also promotes recent fragments into recent.md (--promote auto).
INTERVAL=21600  # 6 hours

while true; do
  sleep "$INTERVAL"
  /usr/local/bin/term-llm memory mine \
    --agent "${AGENT_NAME:-agent}" \
    --promote auto \
    --since 24h \
    --fast-model \
    2>&1
done
