# docker/templates/services/jobs/run
#!/bin/sh
exec /usr/local/bin/term-llm serve jobs \
  --agent "${AGENT_NAME:-agent}" \
  --auth none \
  --yolo \
  2>&1
